# nginx-supervisord-lab/supervisord.conf
[unix_http_server]
file=/run/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

# Tạo hai backend từ cùng một chương trình: port 9000 và 9001.
[program:app]
command=/usr/bin/python3 /opt/lab/app.py --port=90%(process_num)02d --name=worker_%(process_num)02d
process_name=%(program_name)s_%(process_num)02d
numprocs=2
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/app_%(process_num)02d_stdout.log
stderr_logfile=/var/log/supervisor/app_%(process_num)02d_stderr.log
environment=PYTHONUNBUFFERED="1"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=unexpected
stdout_logfile=/var/log/supervisor/nginx_stdout.log
stderr_logfile=/var/log/supervisor/nginx_stderr.log

# Service thêm động được đặt vào thư mục này rồi dùng reread + update.
[include]
files=/etc/supervisor/conf.d/*.conf
